# opencontrail/conf/config/ifmap.ini
[program:ifmap]
directory=/usr/share/irond
command=/usr/bin/java -Dlog4j.configuration=file:////etc/irond/log4j.properties -jar /usr/libexec/irond.jar
priority=440
autostart=true
killasgroup=true
stopsignal=KILL
serverurl=http://localhost:9004
redirect_stderr=true
stdout_logfile=/var/log/contrail/ifmap.log
stderr_logfile=/dev/null
